# launchd/com.openrouter.menu-bar.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openrouter.menubar.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>__REPO_PATH__/collector/openrouter_collector.py</string>
        <string>--verbose</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENROUTER_API_KEY</key>
        <string>__OPENROUTER_API_KEY__</string>
        <key>OPENROUTER_MANAGEMENT_API_KEY</key>
        <string>__OPENROUTER_MANAGEMENT_API_KEY__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>__REPO_PATH__</string>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/collector-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/collector-launchd.log</string>
</dict>
</plist>
